# systemd/batmon.timer
[Unit]
Description=Run batmon every 30 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
